# agent/homelan-agent.service
[Unit]
Description=HomeLan System Monitor Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/homelan/homelan-agent.py
WorkingDirectory=/opt/homelan
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homelan-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/proc /sys
ReadWritePaths=/opt/homelan

[Install]
WantedBy=multi-user.target
